# lifting-lookup-refresh/crontab
0 1 * * * bash /home/ubuntu/lifting-lookup/lifting-lookup-refresh/refresh.sh